# meta-edgefarm-core/recipes-core/hostname2hosts/files/hostname2hosts.service
[Unit]
Description=Sync /etc/hosts with /etc/hostname

[Service]
Type=oneshot
ExecStart=/usr/bin/hostname2hosts.sh

[Install]
WantedBy=multi-user.target
