# meta-edgefarm-core/recipes-overlays/persistent-edgecore/files/persistent-edgecore.service
[Unit]
Description=Create and mount overlay filesystems.
Requires=data.mount
After=data.mount

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/persistent-edgecore.sh

[Install]
WantedBy=multi-user.target
