# meta-edgefarm-core/recipes-containers/docker-prune/files/docker-prunea-at-boot.service
[Unit]
Description=Clean Up Docker Containers
[Service]
ExecStart=/usr/bin/docker container prune --force
[Install]
WantedBy=multi-user.target
