# meta-edgefarm-core/recipes-containers/docker-prune/files/docker-prune-at-boot.service
[Unit]
Description=Clean Up Docker Containers
After=docker.service
Requires=docker.service
Before=kubelet.service

[Service]
ExecStart=/usr/bin/docker container prune --force
[Install]
WantedBy=multi-user.target
